# build/config/common/includes.chroot/etc/sv/alsa/run
#!/bin/sh
# Restore ALSA state and stay running
/usr/sbin/alsactl restore 2>/dev/null || true
exec sleep infinity
