# build/config/common/includes.chroot/etc/sv/busybox-syslog/run
#!/bin/sh
exec /sbin/syslogd -n -O /var/log/syslog 2>&1
